# data/config.conf
[program:messenger_consume_index_site]
    command=php /var/www/symfony/bin/console messenger:consume index_site --time-limit=3600
    process_name=consumer_%(process_num)02d
    max=75
    autostart=true
    autorestart=true
    queue=rabbitmq
    queue_connection=amqp://guest:guest@localhost:5678
    queue_name=messages_low